# cachyos-keyring-sync/system/cachyos-keyring-sync.timer
[Unit]
Description=Weekly CachyOS keyring sync

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
